# tieba-sign/.config/systemd/user/tieba-sign.timer
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=*:00:00

