# torrent-manager/.config/systemd/user/downloading-torrents-resume.timer
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=daily

[Unit]
Description=Start when midnight.
