# torrent-manager/.config/systemd/user/downloading-torrents-pause.service
[Service]
ExecStart=torrent-manager pause
Restart=no
Type=oneshot

[Unit]
Description=Pause downloading torrents
