# updates/.config/systemd/user/flatpak-update.service
[Unit]
After=network-online.target

[Service]
Type=oneshot
ExecStart=flatpak update --no-deploy -y
