# updates/.config/systemd/user/flatpak-update.timer
[Install]
WantedBy=default.target

[Timer]
OnCalendar=daily
