# updates/.config/systemd/user/update@.service
[Unit]
Description=updating for %i
After=network-online.target

[Service]
ExecStart=zsh -c 'update-%i'

