# updates/.config/systemd/user/updates.service
[Unit]
Description=doing some updates
After=network-online.target

[Service]
ExecStart=zsh -c 'updates'
