# updates/.config/systemd/user/updates.timer
[Unit]
Description=doing updates midnight

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
