# updates/.config/systemd/user/update-daily@.timer
[Timer]
OnCalendar=daily
Unit=update@%i.service

[Install]
WantedBy=timers.target

