# updates/.config/systemd/user/yay-update.service
[Unit]
After=network-online.target

[Service]
Type=oneshot
ExecStart=yay -Syuw --noconfirm

